# supervisor.conf
[program:teelebot]
directory=/home/
command=python3 -m teelebot
autostart=true
autorestart=true
startsecs=2
