# supervisor.conf
[program:TeeleBot]
directory=/home/TeeleBot/
command=python3 /home/TeeleBot/teelebot.py
autostart=true
autorestart=true
startsecs=2
