# supervisor.conf
[program:TeeleBot]
directory=/home/
command=python3 -m TeeleBot
autostart=true
autorestart=true
startsecs=2
